# overlay/usr/lib/systemd/system/boomer@.service
[Unit]
Description=Boomer Test Program running on %f
BindsTo=%i
After=%i

[Service]
ExecStart=/usr/bin/boomer -D %f
